# crontab
15 7 * * *	$HOME/MASSDIMM/timDIMM/rotate.cron >> $HOME/cron.log
*/2 * * * *	$HOME/MASSDIMM/timDIMM/mkmovie.sh >> $HOME/cron.log
