# rotate.cron
#!/bin/sh

cd /Users/timdimm/MASSDIMM/timDIMM
./rotate.sh `date -u -v-1d +'%Y/%m%d'`
